# service/epomaker-controller.service
[Unit]
Description="Epomaker keyboard controller service"

[Service]
Type=simple
Restart=no
ExecStart=/bin/bash epomaker-controller.sh
